# dotfiles/crontab
@reboot rm -rf {{HOME}}/tmp && mkdir {{HOME}}/tmp
